# reverse_proxy/crontab
# Auto-renew let's encrypt SSL certificates
0     *     *     *     *      certbot renew && /usr/sbin/nginx -s reload
